# agent-monitor/deploy/agent-monitor.service
[Unit]
Description=Agent Monitor Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/agent-monitor
Environment=AGENT_MONITOR_HOST=127.0.0.1
Environment=AGENT_MONITOR_PORT=8769
ExecStart=/usr/bin/python3 /opt/agent-monitor/run.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
